# scripts/service/displaygrid.service
[Unit]
Description=DisplayGrid Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/displaygrid
ExecStart=/usr/local/bin/node apps/dashboard/.next/standalone/server.js
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=displaygrid
Environment=NODE_ENV=production
EnvironmentFile=/home/pi/displaygrid/.env

[Install]
WantedBy=multi-user.target
